# contrib/macaw.service
[Unit]
Description=Macaw speech-to-text service
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/macaw
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
